# trunk/contrib/rorserver-initscript.in
#!/bin/bash
### BEGIN INIT INFO
# Provides:          rorserver
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Rigs of Rods multiplayer server
# Description:       rorserver is a server for the Rigs of Rods Game
### END INIT INFO

PATH=@CMAKE_INSTALL_PREFIX@:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
NAME=rorserver
DESC="Rigs of Rods Multiplayer Server"

# if you set CRASHREPORTS=1 then the init script will use gdb to catch errors.
# you should use Debug compile mode. Look in /var/log/rorserver/ for the results
CRASHREPORTS=0

USERNAME=rorserver
DAEMON=@CMAKE_INSTALL_PREFIX@/bin/rorserver
CFGDIR=/etc/rorserver
LOGDIR=/var/log/rorserver
PIDDIR=/var/run/rorserver
RESDIR=@CMAKE_INSTALL_PREFIX@/share/rorserver/

[ -x "$DAEMON" ] || exit 0
[ -d "$CFGDIR" ] || exit 0
[ -d "$LOGDIR" ] || exit 0
[ -d "$PIDDIR" ] || exit 0
[ -d "$RESDIR" ] || exit 0

. /lib/lsb/init-functions

case "$1" in
  start)
	CONFIGS=$(ls ${CFGDIR}/*.cfg 2>/dev/null)
	if [[ "${CONFIGS}" == "" ]] ; then
		log_daemon_msg "no configurations found in ${CFGDIR}"
		log_end_msg 1
		exit 0
	fi
	for CONFIG in ${CONFIGS}
	do
		SRVNAME=$(basename ${CONFIG} | awk -F. '{ print $1 }')
		PIDFILE=${PIDDIR}/${SRVNAME}.pid
		LOGFILE=${LOGDIR}/${SRVNAME}.log
		RUNNING=0
		if [[ -f ${PIDFILE} ]] ; then
			PIDNUM=$(cat ${PIDFILE})
			RUNNING=$(ps -p ${PIDNUM} 2>/dev/null | wc -l)
		fi
		if [[ "${RUNNING}" != "2" ]] ; then
			# start
			rm -f ${PIDFILE} >>/dev/null 2>&1
			log_daemon_msg "starting ${NAME} with config ${SRVNAME}"
			echo "LOG STARTED" > ${LOGFILE}
			chown $USERNAME:$USERNAME ${LOGFILE}
			if [[ "${CRASHREPORTS}" == "0" ]]
			then
				# normal startup, no gdb
				start-stop-daemon --start --background --user ${USERNAME} --name ${NAME}-${SRVNAME} --exec ${DAEMON} --make-pidfile --pidfile ${PIDFILE}  -- -c ${CONFIG} -resdir ${RESDIR} -logfilename ${LOGFILE}
				log_end_msg $?
			else
				# gdb startup
				DATENOW=$(date +%s)
				GDB=$(which gdb)
				GDBCMDFILE=${LOGFILE}.${DATENOW}.gdb.cmds
				GDBLOG=${LOGFILE}.${DATENOW}.gdb.log
				echo "set args -c ${CONFIG} -resdir ${RESDIR} -logfilename ${LOGFILE} -logverbosity 0 -fg" > ${GDBCMDFILE}
				echo "run" >> ${GDBCMDFILE}
				echo "set logging file ${GDBLOG}" >> ${GDBCMDFILE}
				echo "set logging on" >> ${GDBCMDFILE}
				echo "bt full" >> ${GDBCMDFILE}
				echo "info frame" >> ${GDBCMDFILE}
				echo "info frame 2" >> ${GDBCMDFILE}
				echo "info frame 3" >> ${GDBCMDFILE}
				echo "quit" >> ${GDBCMDFILE}
"
				start-stop-daemon --start --background --user ${USERNAME} --name ${NAME}-${SRVNAME} --exec ${GDB} --make-pidfile --pidfile ${PIDFILE}  -- ${DAEMON} -x ${GDBCMDFILE}
				log_end_msg $?
			fi
		else
			log_daemon_msg "${NAME} with config ${SRVNAME} already running"
			log_end_msg 0
		fi
	done
	;;
  stop)
	PIDS=$(ls ${PIDDIR}/*.pid 2>/dev/null)
	if [[ "${PIDS}" == "" ]]
	then
		log_daemon_msg "no rorservers detected running"
		log_end_msg 1
		exit 0
	fi
	for PID in ${PIDS}
	do
 		log_daemon_msg "Stopping $NAME $PID"
		start-stop-daemon --stop --pidfile ${PID} --retry=TERM/30/KILL/5
		log_end_msg $?
		#rm ${PID} >>/dev/null 2>&1
	done
	;;
  restart)
	$0 stop && $0 start
	;;
  *)
	echo "Usage: $0 {start|stop|restart}" >&2
	exit 1
	;;
esac

exit 0
